# provisioning/flask/files/myapp.conf
[program:myapp]
command = gunicorn app:app -b localhost:8080
directory = /home/www/
user = www-data
